# jasper.conf
[program:jasper]
command=/home/pi/jasper/run.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/jasper.err.log
stdout_logfile=/var/log/jasper.out.log
